# manifest/armv7l/m/mesa.filelist
/usr/local/include/EGL/eglext_angle.h
/usr/local/include/EGL/eglmesaext.h
/usr/local/include/GL/internal/dri_interface.h
/usr/local/include/gbm.h
/usr/local/include/wayland-protocols/alpha-modifier-v1-enum.h
/usr/local/include/wayland-protocols/commit-timing-v1-enum.h
/usr/local/include/wayland-protocols/content-type-v1-enum.h
/usr/local/include/wayland-protocols/cursor-shape-v1-enum.h
/usr/local/include/wayland-protocols/drm-lease-v1-enum.h
/usr/local/include/wayland-protocols/ext-foreign-toplevel-list-v1-enum.h
/usr/local/include/wayland-protocols/ext-idle-notify-v1-enum.h
/usr/local/include/wayland-protocols/ext-image-capture-source-v1-enum.h
/usr/local/include/wayland-protocols/ext-image-copy-capture-v1-enum.h
/usr/local/include/wayland-protocols/ext-session-lock-v1-enum.h
/usr/local/include/wayland-protocols/ext-transient-seat-v1-enum.h
/usr/local/include/wayland-protocols/fifo-v1-enum.h
/usr/local/include/wayland-protocols/fractional-scale-v1-enum.h
/usr/local/include/wayland-protocols/fullscreen-shell-unstable-v1-enum.h
/usr/local/include/wayland-protocols/idle-inhibit-unstable-v1-enum.h
/usr/local/include/wayland-protocols/input-method-unstable-v1-enum.h
/usr/local/include/wayland-protocols/input-timestamps-unstable-v1-enum.h
/usr/local/include/wayland-protocols/keyboard-shortcuts-inhibit-unstable-v1-enum.h
/usr/local/include/wayland-protocols/linux-dmabuf-unstable-v1-enum.h
/usr/local/include/wayland-protocols/linux-dmabuf-v1-enum.h
/usr/local/include/wayland-protocols/linux-drm-syncobj-v1-enum.h
/usr/local/include/wayland-protocols/linux-explicit-synchronization-unstable-v1-enum.h
/usr/local/include/wayland-protocols/pointer-constraints-unstable-v1-enum.h
/usr/local/include/wayland-protocols/pointer-gestures-unstable-v1-enum.h
/usr/local/include/wayland-protocols/presentation-time-enum.h
/usr/local/include/wayland-protocols/primary-selection-unstable-v1-enum.h
/usr/local/include/wayland-protocols/relative-pointer-unstable-v1-enum.h
/usr/local/include/wayland-protocols/security-context-v1-enum.h
/usr/local/include/wayland-protocols/single-pixel-buffer-v1-enum.h
/usr/local/include/wayland-protocols/tablet-unstable-v1-enum.h
/usr/local/include/wayland-protocols/tablet-unstable-v2-enum.h
/usr/local/include/wayland-protocols/tablet-v2-enum.h
/usr/local/include/wayland-protocols/tearing-control-v1-enum.h
/usr/local/include/wayland-protocols/text-input-unstable-v1-enum.h
/usr/local/include/wayland-protocols/text-input-unstable-v3-enum.h
/usr/local/include/wayland-protocols/viewporter-enum.h
/usr/local/include/wayland-protocols/xdg-activation-v1-enum.h
/usr/local/include/wayland-protocols/xdg-decoration-unstable-v1-enum.h
/usr/local/include/wayland-protocols/xdg-dialog-v1-enum.h
/usr/local/include/wayland-protocols/xdg-foreign-unstable-v1-enum.h
/usr/local/include/wayland-protocols/xdg-foreign-unstable-v2-enum.h
/usr/local/include/wayland-protocols/xdg-output-unstable-v1-enum.h
/usr/local/include/wayland-protocols/xdg-shell-enum.h
/usr/local/include/wayland-protocols/xdg-shell-unstable-v5-enum.h
/usr/local/include/wayland-protocols/xdg-shell-unstable-v6-enum.h
/usr/local/include/wayland-protocols/xdg-system-bell-v1-enum.h
/usr/local/include/wayland-protocols/xdg-toplevel-drag-v1-enum.h
/usr/local/include/wayland-protocols/xdg-toplevel-icon-v1-enum.h
/usr/local/include/wayland-protocols/xwayland-keyboard-grab-unstable-v1-enum.h
/usr/local/include/wayland-protocols/xwayland-shell-v1-enum.h
/usr/local/include/xa_composite.h
/usr/local/include/xa_context.h
/usr/local/include/xa_tracker.h
/usr/local/lib/dri/armada-drm_dri.so
/usr/local/lib/dri/etnaviv_dri.so
/usr/local/lib/dri/exynos_dri.so
/usr/local/lib/dri/gm12u320_dri.so
/usr/local/lib/dri/hdlcd_dri.so
/usr/local/lib/dri/hx8357d_dri.so
/usr/local/lib/dri/ili9163_dri.so
/usr/local/lib/dri/ili9225_dri.so
/usr/local/lib/dri/ili9341_dri.so
/usr/local/lib/dri/ili9486_dri.so
/usr/local/lib/dri/imx-dcss_dri.so
/usr/local/lib/dri/imx-drm_dri.so
/usr/local/lib/dri/imx-lcdif_dri.so
/usr/local/lib/dri/ingenic-drm_dri.so
/usr/local/lib/dri/iris_dri.so
/usr/local/lib/dri/kgsl_dri.so
/usr/local/lib/dri/kirin_dri.so
/usr/local/lib/dri/kms_swrast_dri.so
/usr/local/lib/dri/komeda_dri.so
/usr/local/lib/dri/libdril_dri.so
/usr/local/lib/dri/lima_dri.so
/usr/local/lib/dri/mali-dp_dri.so
/usr/local/lib/dri/mcde_dri.so
/usr/local/lib/dri/mediatek_dri.so
/usr/local/lib/dri/meson_dri.so
/usr/local/lib/dri/mi0283qt_dri.so
/usr/local/lib/dri/msm_dri.so
/usr/local/lib/dri/mxsfb-drm_dri.so
/usr/local/lib/dri/nouveau_dri.so
/usr/local/lib/dri/nouveau_drv_video.so
/usr/local/lib/dri/panel-mipi-dbi_dri.so
/usr/local/lib/dri/panfrost_dri.so
/usr/local/lib/dri/panthor_dri.so
/usr/local/lib/dri/pl111_dri.so
/usr/local/lib/dri/rcar-du_dri.so
/usr/local/lib/dri/repaper_dri.so
/usr/local/lib/dri/rockchip_dri.so
/usr/local/lib/dri/rzg2l-du_dri.so
/usr/local/lib/dri/ssd130x_dri.so
/usr/local/lib/dri/st7586_dri.so
/usr/local/lib/dri/st7735r_dri.so
/usr/local/lib/dri/sti_dri.so
/usr/local/lib/dri/stm_dri.so
/usr/local/lib/dri/sun4i-drm_dri.so
/usr/local/lib/dri/swrast_dri.so
/usr/local/lib/dri/tegra_dri.so
/usr/local/lib/dri/tls
/usr/local/lib/dri/udl_dri.so
/usr/local/lib/dri/v3d_dri.so
/usr/local/lib/dri/virtio_gpu_dri.so
/usr/local/lib/dri/virtio_gpu_drv_video.so
/usr/local/lib/dri/vkms_dri.so
/usr/local/lib/dri/vmwgfx_dri.so
/usr/local/lib/dri/zink_dri.so
/usr/local/lib/dri/zynqmp-dpsub_dri.so
/usr/local/lib/gbm/dri_gbm.so
/usr/local/lib/gbm/tls/pvr_gbm.so
/usr/local/lib/libEGL_mesa.so
/usr/local/lib/libEGL_mesa.so.0
/usr/local/lib/libEGL_mesa.so.0.0.0
/usr/local/lib/libGLX_mesa.so
/usr/local/lib/libGLX_mesa.so.0
/usr/local/lib/libGLX_mesa.so.0.0.0
/usr/local/lib/libgallium-24.3.0-rc2.so
/usr/local/lib/libgbm.so
/usr/local/lib/libgbm.so.1
/usr/local/lib/libgbm.so.1.0.0
/usr/local/lib/libglapi.so
/usr/local/lib/libglapi.so.0
/usr/local/lib/libglapi.so.0.0.0
/usr/local/lib/libvulkan_intel.so
/usr/local/lib/libvulkan_lvp.so
/usr/local/lib/libxatracker.so
/usr/local/lib/libxatracker.so.2
/usr/local/lib/libxatracker.so.2.5.0
/usr/local/lib/pkgconfig/dri.pc
/usr/local/lib/pkgconfig/gbm.pc
/usr/local/lib/pkgconfig/xatracker.pc
/usr/local/lib/vdpau/libvdpau_nouveau.so
/usr/local/lib/vdpau/libvdpau_nouveau.so.1
/usr/local/lib/vdpau/libvdpau_nouveau.so.1.0
/usr/local/lib/vdpau/libvdpau_nouveau.so.1.0.0
/usr/local/lib/vdpau/libvdpau_virtio_gpu.so
/usr/local/lib/vdpau/libvdpau_virtio_gpu.so.1
/usr/local/lib/vdpau/libvdpau_virtio_gpu.so.1.0
/usr/local/lib/vdpau/libvdpau_virtio_gpu.so.1.0.0
/usr/local/share/drirc.d/00-mesa-defaults.conf
/usr/local/share/glvnd/egl_vendor.d/50_mesa.json
/usr/local/share/pkgconfig/wayland-protocols.pc
/usr/local/share/vulkan/icd.d/intel_icd.armv8l.json
/usr/local/share/vulkan/icd.d/lvp_icd.armv8l.json
/usr/local/share/wayland-protocols/stable/linux-dmabuf/linux-dmabuf-v1.xml
/usr/local/share/wayland-protocols/stable/presentation-time/presentation-time.xml
/usr/local/share/wayland-protocols/stable/tablet/tablet-v2.xml
/usr/local/share/wayland-protocols/stable/viewporter/viewporter.xml
/usr/local/share/wayland-protocols/stable/xdg-shell/xdg-shell.xml
/usr/local/share/wayland-protocols/staging/alpha-modifier/alpha-modifier-v1.xml
/usr/local/share/wayland-protocols/staging/commit-timing/commit-timing-v1.xml
/usr/local/share/wayland-protocols/staging/content-type/content-type-v1.xml
/usr/local/share/wayland-protocols/staging/cursor-shape/cursor-shape-v1.xml
/usr/local/share/wayland-protocols/staging/drm-lease/drm-lease-v1.xml
/usr/local/share/wayland-protocols/staging/ext-foreign-toplevel-list/ext-foreign-toplevel-list-v1.xml
/usr/local/share/wayland-protocols/staging/ext-idle-notify/ext-idle-notify-v1.xml
/usr/local/share/wayland-protocols/staging/ext-image-capture-source/ext-image-capture-source-v1.xml
/usr/local/share/wayland-protocols/staging/ext-image-copy-capture/ext-image-copy-capture-v1.xml
/usr/local/share/wayland-protocols/staging/ext-session-lock/ext-session-lock-v1.xml
/usr/local/share/wayland-protocols/staging/ext-transient-seat/ext-transient-seat-v1.xml
/usr/local/share/wayland-protocols/staging/fifo/fifo-v1.xml
/usr/local/share/wayland-protocols/staging/fractional-scale/fractional-scale-v1.xml
/usr/local/share/wayland-protocols/staging/linux-drm-syncobj/linux-drm-syncobj-v1.xml
/usr/local/share/wayland-protocols/staging/security-context/security-context-v1.xml
/usr/local/share/wayland-protocols/staging/single-pixel-buffer/single-pixel-buffer-v1.xml
/usr/local/share/wayland-protocols/staging/tearing-control/tearing-control-v1.xml
/usr/local/share/wayland-protocols/staging/xdg-activation/xdg-activation-v1.xml
/usr/local/share/wayland-protocols/staging/xdg-dialog/xdg-dialog-v1.xml
/usr/local/share/wayland-protocols/staging/xdg-system-bell/xdg-system-bell-v1.xml
/usr/local/share/wayland-protocols/staging/xdg-toplevel-drag/xdg-toplevel-drag-v1.xml
/usr/local/share/wayland-protocols/staging/xdg-toplevel-icon/xdg-toplevel-icon-v1.xml
/usr/local/share/wayland-protocols/staging/xwayland-shell/xwayland-shell-v1.xml
/usr/local/share/wayland-protocols/unstable/fullscreen-shell/fullscreen-shell-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/idle-inhibit/idle-inhibit-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/input-method/input-method-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/input-timestamps/input-timestamps-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/keyboard-shortcuts-inhibit/keyboard-shortcuts-inhibit-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/linux-dmabuf/linux-dmabuf-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/linux-explicit-synchronization/linux-explicit-synchronization-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/pointer-constraints/pointer-constraints-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/pointer-gestures/pointer-gestures-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/primary-selection/primary-selection-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/relative-pointer/relative-pointer-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/tablet/tablet-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/tablet/tablet-unstable-v2.xml
/usr/local/share/wayland-protocols/unstable/text-input/text-input-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/text-input/text-input-unstable-v3.xml
/usr/local/share/wayland-protocols/unstable/xdg-decoration/xdg-decoration-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/xdg-foreign/xdg-foreign-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/xdg-foreign/xdg-foreign-unstable-v2.xml
/usr/local/share/wayland-protocols/unstable/xdg-output/xdg-output-unstable-v1.xml
/usr/local/share/wayland-protocols/unstable/xdg-shell/xdg-shell-unstable-v5.xml
/usr/local/share/wayland-protocols/unstable/xdg-shell/xdg-shell-unstable-v6.xml
/usr/local/share/wayland-protocols/unstable/xwayland-keyboard-grab/xwayland-keyboard-grab-unstable-v1.xml
