# manifest/x86_64/v/vdev.filelist
/usr/local/etc/init.d/vdev
/usr/local/etc/vdev/actions/001-sysfs-block.act
/usr/local/etc/vdev/actions/001-sysfs-char.act
/usr/local/etc/vdev/actions/002-modalias.act
/usr/local/etc/vdev/actions/003-firmware.act
/usr/local/etc/vdev/actions/agpgart.act
/usr/local/etc/vdev/actions/block.act
/usr/local/etc/vdev/actions/btrfs-control.act
/usr/local/etc/vdev/actions/char.act
/usr/local/etc/vdev/actions/cpu-msr.act
/usr/local/etc/vdev/actions/cpu.act
/usr/local/etc/vdev/actions/cuse.act
/usr/local/etc/vdev/actions/device-mapper.act
/usr/local/etc/vdev/actions/dialout.act
/usr/local/etc/vdev/actions/disk-part.act
/usr/local/etc/vdev/actions/disk.act
/usr/local/etc/vdev/actions/dm-disk.act
/usr/local/etc/vdev/actions/dri.act
/usr/local/etc/vdev/actions/err.act
/usr/local/etc/vdev/actions/fb.act
/usr/local/etc/vdev/actions/full.act
/usr/local/etc/vdev/actions/fuse-mounts.act
/usr/local/etc/vdev/actions/fuse.act
/usr/local/etc/vdev/actions/hwdb.act
/usr/local/etc/vdev/actions/ifname.act
/usr/local/etc/vdev/actions/input.act
/usr/local/etc/vdev/actions/irlpt.act
/usr/local/etc/vdev/actions/kmsg.act
/usr/local/etc/vdev/actions/kvm.act
/usr/local/etc/vdev/actions/legousbtower.act
/usr/local/etc/vdev/actions/lircN.act
/usr/local/etc/vdev/actions/loop-control.act
/usr/local/etc/vdev/actions/loop.act
/usr/local/etc/vdev/actions/lp.act
/usr/local/etc/vdev/actions/mISDNtimer.act
/usr/local/etc/vdev/actions/mem.act
/usr/local/etc/vdev/actions/mmtimer.act
/usr/local/etc/vdev/actions/mwave.act
/usr/local/etc/vdev/actions/net-tun.act
/usr/local/etc/vdev/actions/netdev.act
/usr/local/etc/vdev/actions/null.act
/usr/local/etc/vdev/actions/nvram.act
/usr/local/etc/vdev/actions/optical.act
/usr/local/etc/vdev/actions/parport.act
/usr/local/etc/vdev/actions/pktcdvd.act
/usr/local/etc/vdev/actions/pktcdvdN.act
/usr/local/etc/vdev/actions/port.act
/usr/local/etc/vdev/actions/power-switch.act
/usr/local/etc/vdev/actions/printer.act
/usr/local/etc/vdev/actions/ptmx.act
/usr/local/etc/vdev/actions/qftN.act
/usr/local/etc/vdev/actions/random.act
/usr/local/etc/vdev/actions/raw.act
/usr/local/etc/vdev/actions/rawctl.act
/usr/local/etc/vdev/actions/rfkill.act
/usr/local/etc/vdev/actions/rtc.act
/usr/local/etc/vdev/actions/sch.act
/usr/local/etc/vdev/actions/sclp_line.act
/usr/local/etc/vdev/actions/sgi.act
/usr/local/etc/vdev/actions/sonypi.act
/usr/local/etc/vdev/actions/sound-alsa-control.act
/usr/local/etc/vdev/actions/sound-audio.act
/usr/local/etc/vdev/actions/tty.act
/usr/local/etc/vdev/actions/tty3270.act
/usr/local/etc/vdev/actions/ttyN.act
/usr/local/etc/vdev/actions/ttyS.act
/usr/local/etc/vdev/actions/ttysclp.act
/usr/local/etc/vdev/actions/usb-bus.act
/usr/local/etc/vdev/actions/v4l.act
/usr/local/etc/vdev/actions/vcs.act
/usr/local/etc/vdev/actions/virtio.act
/usr/local/etc/vdev/actions/virtualbox-usb.act
/usr/local/etc/vdev/actions/virtualbox-usb_device.act
/usr/local/etc/vdev/actions/z90crypt.act
/usr/local/etc/vdev/actions/zzz-event-push.act
/usr/local/etc/vdev/actions/zzz-udev-compat.act
/usr/local/etc/vdev/ifnames.conf
/usr/local/etc/vdev/vdevd.conf
/usr/local/include/libudev.h
/usr/local/lib64/libudev.so.1
/usr/local/lib64/libudev.so.1.5.2
/usr/local/lib64/vdev/VBoxCreateUSBNode.sh
/usr/local/lib64/vdev/block.sh
/usr/local/lib64/vdev/char.sh
/usr/local/lib64/vdev/daemonlet
/usr/local/lib64/vdev/dev-setup.sh
/usr/local/lib64/vdev/disk.sh
/usr/local/lib64/vdev/dm-disk.sh
/usr/local/lib64/vdev/echo_n
/usr/local/lib64/vdev/event-push.sh
/usr/local/lib64/vdev/event-put
/usr/local/lib64/vdev/firmware.sh
/usr/local/lib64/vdev/gen-ifnames.sh
/usr/local/lib64/vdev/hwdb-props.sh
/usr/local/lib64/vdev/hwdb.sh
/usr/local/lib64/vdev/hwdb/hwdb.squashfs
/usr/local/lib64/vdev/ifname.sh
/usr/local/lib64/vdev/input.sh
/usr/local/lib64/vdev/modprobe.sh
/usr/local/lib64/vdev/netdev.sh
/usr/local/lib64/vdev/optical.sh
/usr/local/lib64/vdev/permissions.sh
/usr/local/lib64/vdev/power-switch.sh
/usr/local/lib64/vdev/sound.sh
/usr/local/lib64/vdev/stat_ata
/usr/local/lib64/vdev/stat_bus
/usr/local/lib64/vdev/stat_input
/usr/local/lib64/vdev/stat_net
/usr/local/lib64/vdev/stat_optical
/usr/local/lib64/vdev/stat_path
/usr/local/lib64/vdev/stat_scsi
/usr/local/lib64/vdev/stat_usb
/usr/local/lib64/vdev/stat_v4l
/usr/local/lib64/vdev/subr-event.sh
/usr/local/lib64/vdev/subr-hwdb.sh
/usr/local/lib64/vdev/subr.sh
/usr/local/lib64/vdev/sysfs.sh
/usr/local/lib64/vdev/udev-compat.sh
/usr/local/lib64/vdev/usb-name.sh
/usr/local/lib64/vdev/v4l.sh
/usr/local/sbin/vdevd
/usr/local/sbin/vdevfs
/usr/local/share/pkgconfig/libudev.pc
