# manifest/i686/v/v2ray.filelist
/usr/local/bin/v2ray
/usr/local/etc/v2ray/config.json
/usr/local/share/v2ray/config.json
/usr/local/share/v2ray/geoip-only-cn-private.dat
/usr/local/share/v2ray/geoip.dat
/usr/local/share/v2ray/geosite.dat
/usr/local/share/v2ray/systemd/system/v2ray.service
/usr/local/share/v2ray/systemd/system/v2ray@.service
/usr/local/share/v2ray/v2ray
/usr/local/share/v2ray/vpoint_socks_vmess.json
/usr/local/share/v2ray/vpoint_vmess_freedom.json
