# manifest/armv7l/i/iputils.filelist
/usr/local/sbin/rdisc
/usr/local/sbin/ninfod
/usr/local/etc/init.d/ninfod.sh
/usr/local/share/locale/de/LC_MESSAGES/iputils.mo
/usr/local/share/locale/cs/LC_MESSAGES/iputils.mo
/usr/local/share/locale/uk/LC_MESSAGES/iputils.mo
/usr/local/share/locale/tr/LC_MESSAGES/iputils.mo
/usr/local/share/locale/pt_BR/LC_MESSAGES/iputils.mo
/usr/local/share/locale/fr/LC_MESSAGES/iputils.mo
/usr/local/share/locale/fi/LC_MESSAGES/iputils.mo
/usr/local/share/locale/zh_CN/LC_MESSAGES/iputils.mo
/usr/local/share/locale/ko/LC_MESSAGES/iputils.mo
/usr/local/share/locale/id/LC_MESSAGES/iputils.mo
/usr/local/share/locale/ja/LC_MESSAGES/iputils.mo
/usr/local/bin/arping
/usr/local/bin/tracepath
/usr/local/bin/ping
/usr/local/bin/clockdiff
