# manifest/x86_64/m/mesa.filelist
/usr/local/include/EGL/eglext_angle.h
/usr/local/include/EGL/eglmesaext.h
/usr/local/include/gbm.h
/usr/local/include/GL/internal/dri_interface.h
/usr/local/include/xa_composite.h
/usr/local/include/xa_context.h
/usr/local/include/xa_tracker.h
/usr/local/lib64/dri/crocus_dri.so
/usr/local/lib64/dri/i915_dri.so
/usr/local/lib64/dri/iris_dri.so
/usr/local/lib64/dri/kms_swrast_dri.so
/usr/local/lib64/dri/nouveau_dri.so
/usr/local/lib64/dri/nouveau_drv_video.so
/usr/local/lib64/dri/r300_dri.so
/usr/local/lib64/dri/r600_dri.so
/usr/local/lib64/dri/r600_drv_video.so
/usr/local/lib64/dri/radeonsi_dri.so
/usr/local/lib64/dri/radeonsi_drv_video.so
/usr/local/lib64/dri/swrast_dri.so
/usr/local/lib64/dri/tls
/usr/local/lib64/dri/virtio_gpu_dri.so
/usr/local/lib64/dri/virtio_gpu_drv_video.so
/usr/local/lib64/dri/vmwgfx_dri.so
/usr/local/lib64/dri/zink_dri.so
/usr/local/lib64/gbm/tls/amdgpu_gbm.so
/usr/local/lib64/gbm/tls/i915_gbm.so
/usr/local/lib64/libEGL_mesa.so
/usr/local/lib64/libEGL_mesa.so.0
/usr/local/lib64/libEGL_mesa.so.0.0.0
/usr/local/lib64/libgbm.so
/usr/local/lib64/libgbm.so.1
/usr/local/lib64/libgbm.so.1.0.0
/usr/local/lib64/libglapi.so
/usr/local/lib64/libglapi.so.0
/usr/local/lib64/libglapi.so.0.0.0
/usr/local/lib64/libGLX_mesa.so
/usr/local/lib64/libGLX_mesa.so.0
/usr/local/lib64/libGLX_mesa.so.0.0.0
/usr/local/lib64/libvulkan_intel_hasvk.so
/usr/local/lib64/libvulkan_intel.so
/usr/local/lib64/libvulkan_lvp.so
/usr/local/lib64/libvulkan_radeon.so
/usr/local/lib64/libxatracker.so
/usr/local/lib64/libxatracker.so.2
/usr/local/lib64/libxatracker.so.2.5.0
/usr/local/lib64/pkgconfig/dri.pc
/usr/local/lib64/pkgconfig/gbm.pc
/usr/local/lib64/pkgconfig/xatracker.pc
/usr/local/lib64/vdpau/libvdpau_nouveau.so
/usr/local/lib64/vdpau/libvdpau_nouveau.so.1
/usr/local/lib64/vdpau/libvdpau_nouveau.so.1.0
/usr/local/lib64/vdpau/libvdpau_nouveau.so.1.0.0
/usr/local/lib64/vdpau/libvdpau_r300.so
/usr/local/lib64/vdpau/libvdpau_r300.so.1
/usr/local/lib64/vdpau/libvdpau_r300.so.1.0
/usr/local/lib64/vdpau/libvdpau_r300.so.1.0.0
/usr/local/lib64/vdpau/libvdpau_r600.so
/usr/local/lib64/vdpau/libvdpau_r600.so.1
/usr/local/lib64/vdpau/libvdpau_r600.so.1.0
/usr/local/lib64/vdpau/libvdpau_r600.so.1.0.0
/usr/local/lib64/vdpau/libvdpau_radeonsi.so
/usr/local/lib64/vdpau/libvdpau_radeonsi.so.1
/usr/local/lib64/vdpau/libvdpau_radeonsi.so.1.0
/usr/local/lib64/vdpau/libvdpau_radeonsi.so.1.0.0
/usr/local/lib64/vdpau/libvdpau_virtio_gpu.so
/usr/local/lib64/vdpau/libvdpau_virtio_gpu.so.1
/usr/local/lib64/vdpau/libvdpau_virtio_gpu.so.1.0
/usr/local/lib64/vdpau/libvdpau_virtio_gpu.so.1.0.0
/usr/local/share/drirc.d/00-mesa-defaults.conf
/usr/local/share/drirc.d/00-radv-defaults.conf
/usr/local/share/glvnd/egl_vendor.d/50_mesa.json
/usr/local/share/vulkan/icd.d/intel_hasvk_icd.x86_64.json
/usr/local/share/vulkan/icd.d/intel_icd.x86_64.json
/usr/local/share/vulkan/icd.d/lvp_icd.x86_64.json
/usr/local/share/vulkan/icd.d/radeon_icd.x86_64.json
