# manifest/x86_64/m/mesa.filelist
# Total size: 133486151
/usr/local/include/EGL/eglext_angle.h
/usr/local/include/EGL/eglmesaext.h
/usr/local/include/GL/internal/dri_interface.h
/usr/local/include/gbm.h
/usr/local/include/gbm_backend_abi.h
/usr/local/lib64/dri/apple_dri.so
/usr/local/lib64/dri/armada-drm_dri.so
/usr/local/lib64/dri/crocus_dri.so
/usr/local/lib64/dri/exynos_dri.so
/usr/local/lib64/dri/gm12u320_dri.so
/usr/local/lib64/dri/hdlcd_dri.so
/usr/local/lib64/dri/hx8357d_dri.so
/usr/local/lib64/dri/i915_dri.so
/usr/local/lib64/dri/ili9163_dri.so
/usr/local/lib64/dri/ili9225_dri.so
/usr/local/lib64/dri/ili9341_dri.so
/usr/local/lib64/dri/ili9486_dri.so
/usr/local/lib64/dri/imx-dcss_dri.so
/usr/local/lib64/dri/imx-drm_dri.so
/usr/local/lib64/dri/imx-lcdif_dri.so
/usr/local/lib64/dri/ingenic-drm_dri.so
/usr/local/lib64/dri/iris_dri.so
/usr/local/lib64/dri/kirin_dri.so
/usr/local/lib64/dri/kms_swrast_dri.so
/usr/local/lib64/dri/komeda_dri.so
/usr/local/lib64/dri/libdril_dri.so
/usr/local/lib64/dri/mali-dp_dri.so
/usr/local/lib64/dri/mcde_dri.so
/usr/local/lib64/dri/mediatek_dri.so
/usr/local/lib64/dri/meson_dri.so
/usr/local/lib64/dri/mi0283qt_dri.so
/usr/local/lib64/dri/mxsfb-drm_dri.so
/usr/local/lib64/dri/nouveau_dri.so
/usr/local/lib64/dri/nouveau_drv_video.so
/usr/local/lib64/dri/panel-mipi-dbi_dri.so
/usr/local/lib64/dri/pl111_dri.so
/usr/local/lib64/dri/r300_dri.so
/usr/local/lib64/dri/r600_dri.so
/usr/local/lib64/dri/r600_drv_video.so
/usr/local/lib64/dri/radeonsi_dri.so
/usr/local/lib64/dri/radeonsi_drv_video.so
/usr/local/lib64/dri/rcar-du_dri.so
/usr/local/lib64/dri/repaper_dri.so
/usr/local/lib64/dri/rockchip_dri.so
/usr/local/lib64/dri/rzg2l-du_dri.so
/usr/local/lib64/dri/ssd130x_dri.so
/usr/local/lib64/dri/st7586_dri.so
/usr/local/lib64/dri/st7735r_dri.so
/usr/local/lib64/dri/sti_dri.so
/usr/local/lib64/dri/stm_dri.so
/usr/local/lib64/dri/sun4i-drm_dri.so
/usr/local/lib64/dri/swrast_dri.so
/usr/local/lib64/dri/tls
/usr/local/lib64/dri/tls/apple_dri.so
/usr/local/lib64/dri/tls/armada-drm_dri.so
/usr/local/lib64/dri/tls/crocus_dri.so
/usr/local/lib64/dri/tls/exynos_dri.so
/usr/local/lib64/dri/tls/gm12u320_dri.so
/usr/local/lib64/dri/tls/hdlcd_dri.so
/usr/local/lib64/dri/tls/hx8357d_dri.so
/usr/local/lib64/dri/tls/i915_dri.so
/usr/local/lib64/dri/tls/ili9163_dri.so
/usr/local/lib64/dri/tls/ili9225_dri.so
/usr/local/lib64/dri/tls/ili9341_dri.so
/usr/local/lib64/dri/tls/ili9486_dri.so
/usr/local/lib64/dri/tls/imx-dcss_dri.so
/usr/local/lib64/dri/tls/imx-drm_dri.so
/usr/local/lib64/dri/tls/imx-lcdif_dri.so
/usr/local/lib64/dri/tls/ingenic-drm_dri.so
/usr/local/lib64/dri/tls/iris_dri.so
/usr/local/lib64/dri/tls/kirin_dri.so
/usr/local/lib64/dri/tls/kms_swrast_dri.so
/usr/local/lib64/dri/tls/komeda_dri.so
/usr/local/lib64/dri/tls/libdril_dri.so
/usr/local/lib64/dri/tls/mali-dp_dri.so
/usr/local/lib64/dri/tls/mcde_dri.so
/usr/local/lib64/dri/tls/mediatek_dri.so
/usr/local/lib64/dri/tls/meson_dri.so
/usr/local/lib64/dri/tls/mi0283qt_dri.so
/usr/local/lib64/dri/tls/mxsfb-drm_dri.so
/usr/local/lib64/dri/tls/nouveau_dri.so
/usr/local/lib64/dri/tls/nouveau_drv_video.so
/usr/local/lib64/dri/tls/panel-mipi-dbi_dri.so
/usr/local/lib64/dri/tls/pl111_dri.so
/usr/local/lib64/dri/tls/r300_dri.so
/usr/local/lib64/dri/tls/r600_dri.so
/usr/local/lib64/dri/tls/r600_drv_video.so
/usr/local/lib64/dri/tls/radeonsi_dri.so
/usr/local/lib64/dri/tls/radeonsi_drv_video.so
/usr/local/lib64/dri/tls/rcar-du_dri.so
/usr/local/lib64/dri/tls/repaper_dri.so
/usr/local/lib64/dri/tls/rockchip_dri.so
/usr/local/lib64/dri/tls/rzg2l-du_dri.so
/usr/local/lib64/dri/tls/ssd130x_dri.so
/usr/local/lib64/dri/tls/st7586_dri.so
/usr/local/lib64/dri/tls/st7735r_dri.so
/usr/local/lib64/dri/tls/sti_dri.so
/usr/local/lib64/dri/tls/stm_dri.so
/usr/local/lib64/dri/tls/sun4i-drm_dri.so
/usr/local/lib64/dri/tls/swrast_dri.so
/usr/local/lib64/dri/tls/tls
/usr/local/lib64/dri/tls/udl_dri.so
/usr/local/lib64/dri/tls/virtio_gpu_dri.so
/usr/local/lib64/dri/tls/virtio_gpu_drv_video.so
/usr/local/lib64/dri/tls/vkms_dri.so
/usr/local/lib64/dri/tls/vmwgfx_dri.so
/usr/local/lib64/dri/tls/zink_dri.so
/usr/local/lib64/dri/tls/zynqmp-dpsub_dri.so
/usr/local/lib64/dri/udl_dri.so
/usr/local/lib64/dri/virtio_gpu_dri.so
/usr/local/lib64/dri/virtio_gpu_drv_video.so
/usr/local/lib64/dri/vkms_dri.so
/usr/local/lib64/dri/vmwgfx_dri.so
/usr/local/lib64/dri/zink_dri.so
/usr/local/lib64/dri/zynqmp-dpsub_dri.so
/usr/local/lib64/gbm/dri_gbm.so
/usr/local/lib64/gbm/tls/amdgpu_gbm.so
/usr/local/lib64/gbm/tls/i915_gbm.so
/usr/local/lib64/libEGL_mesa.so
/usr/local/lib64/libEGL_mesa.so.0
/usr/local/lib64/libEGL_mesa.so.0.0.0
/usr/local/lib64/libGLX_mesa.so
/usr/local/lib64/libGLX_mesa.so.0
/usr/local/lib64/libGLX_mesa.so.0.0.0
/usr/local/lib64/libgallium-26.0.0.so
/usr/local/lib64/libgbm.so
/usr/local/lib64/libgbm.so.1
/usr/local/lib64/libgbm.so.1.0.0
/usr/local/lib64/libvulkan_intel.so
/usr/local/lib64/libvulkan_intel_hasvk.so
/usr/local/lib64/libvulkan_lvp.so
/usr/local/lib64/libvulkan_radeon.so
/usr/local/lib64/pkgconfig/dri.pc
/usr/local/lib64/pkgconfig/gbm.pc
/usr/local/share/drirc.d/00-mesa-defaults.conf
/usr/local/share/drirc.d/00-radv-defaults.conf
/usr/local/share/glvnd/egl_vendor.d/50_mesa.json
/usr/local/share/vulkan/icd.d/intel_hasvk_icd.x86_64.json
/usr/local/share/vulkan/icd.d/intel_icd.x86_64.json
/usr/local/share/vulkan/icd.d/lvp_icd.x86_64.json
/usr/local/share/vulkan/icd.d/radeon_icd.x86_64.json
