# manifest/c/consolekit.filelist
/usr/local/usr/bin/console-kit-daemon
/usr/local/usr/bin/ck-log-system-start
/usr/local/usr/bin/ck-log-system-restart
/usr/local/usr/bin/ck-log-system-stop
/usr/local/etc/ConsoleKit/seats.d/00-primary.seat
/usr/local/etc/logrotate.d/consolekit
/usr/local/etc/X11/xinit/xinitrc.d/90-consolekit
/usr/local/etc/dbus-1/system.d/ConsoleKit.conf
/usr/local/lib64/ConsoleKit/scripts/ck-system-stop
/usr/local/lib64/ConsoleKit/scripts/ck-system-restart
/usr/local/lib64/ConsoleKit/scripts/ck-system-suspend
/usr/local/lib64/ConsoleKit/scripts/ck-system-hibernate
/usr/local/lib64/ConsoleKit/scripts/ck-system-hybridsleep
/usr/local/lib64/ConsoleKit/ck-collect-session-info
/usr/local/lib64/ConsoleKit/ck-get-x11-server-pid
/usr/local/lib64/ConsoleKit/ck-get-x11-display-device
/usr/local/lib64/ConsoleKit/ck-remove-directory
/usr/local/lib64/ConsoleKit/udev-acl
/usr/local/lib64/libconsolekit.so.1.0.0
/usr/local/lib64/libconsolekit.la
/usr/local/lib64/pkgconfig/libconsolekit.pc
/usr/local/lib64/pkgconfig/ck-connector.pc
/usr/local/lib64/girepository-1.0/libconsolekit-1.0.typelib
/usr/local/lib64/libck-connector.so.0.0.0
/usr/local/lib64/libck-connector.la
/usr/local/lib64/security/pam_ck_connector.so
/usr/local/lib64/security/pam_ck_connector.la
/usr/local/lib64/tmpfiles.d/consolekit.conf
/usr/local/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Manager.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Seat.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Session.xml
/usr/local/share/dbus-1/system-services/org.freedesktop.ConsoleKit.service
/usr/local/share/gir-1.0/libconsolekit-1.0.gir
/usr/local/share/polkit-1/actions/org.freedesktop.consolekit.policy
/usr/local/share/man/man1/ck-history.1.gz
/usr/local/share/man/man1/ck-launch-session.1.gz
/usr/local/share/man/man1/ck-list-sessions.1.gz
/usr/local/share/man/man1/console-kit-daemon.1m.gz
/usr/local/share/man/man8/pam_ck_connector.8.gz
/usr/local/share/locale/bg/LC_MESSAGES/ConsoleKit2.mo
/usr/local/share/locale/es/LC_MESSAGES/ConsoleKit2.mo
/usr/local/share/locale/es_419/LC_MESSAGES/ConsoleKit2.mo
/usr/local/share/locale/ru/LC_MESSAGES/ConsoleKit2.mo
/usr/local/include/ConsoleKit/libconsolekit.h
/usr/local/include/ConsoleKit/ck-connector/ck-connector.h
/usr/local/bin/ck-launch-session
/usr/local/bin/ck-list-sessions
/usr/local/bin/ck-history
/usr/local/libexec/rules.d/70-udev-acl.rules
/usr/local/libexec/rules.d/71-udev-seat.rules
/usr/local/lib64/ConsoleKit/run-seat.d/udev-acl.ck
/usr/local/lib64/libconsolekit.so.1
/usr/local/lib64/libconsolekit.so
/usr/local/lib64/libck-connector.so.0
/usr/local/lib64/libck-connector.so
/usr/local/libexec/udev-acl
