# manifest/armv7l/c/chrx.filelist
/usr/local/bin/chrx
/usr/local/bin/chrx-install
/usr/local/bin/chrx-install-chroot
/usr/local/bin/chrx-install-extras
/usr/local/bin/chrx-setup-storage
/usr/local/etc/chrx-devices
/usr/local/etc/chrx-files/etc/X11/xorg.conf
/usr/local/etc/chrx-files/etc/X11/xorg.conf.d/20-intel-graphics.conf
/usr/local/etc/chrx-files/etc/X11/xorg.conf.d/50-cyapa-touchpad.conf
/usr/local/etc/chrx-files/etc/X11/xorg.conf.d/70-synaptics.conf
/usr/local/etc/chrx-files/etc/dracut.conf.d/chrx.conf
/usr/local/etc/chrx-files/etc/dracut.conf.d/mmc.conf
/usr/local/etc/chrx-files/etc/init.d/chrx-firstboot
/usr/local/etc/chrx-files/etc/initramfs-tools/scripts/init-top/unbind_ehci
/usr/local/etc/chrx-files/etc/modprobe.d/ath9k.conf
/usr/local/etc/chrx-files/etc/modprobe.d/snd-hda-intel.conf
/usr/local/etc/chrx-files/etc/pm/sleep.d/50-audio-intel-driver
/usr/local/etc/chrx-files/etc/pm/sleep.d/89-wifi-ath9k-module
/usr/local/etc/chrx-files/etc/rc.local
/usr/local/etc/chrx-files/etc/selinux/config
/usr/local/etc/chrx-files/etc/skel/.config/autostart/chrx_touchpad_support.desktop
/usr/local/etc/chrx-files/etc/skel/.config/autostart/xbindkeys.desktop
/usr/local/etc/chrx-files/etc/skel/.config/lxpanel/Lubuntu/config
/usr/local/etc/chrx-files/etc/skel/.config/lxpanel/Lubuntu/panels/panel
/usr/local/etc/chrx-files/etc/skel/.local/bin/chrx_touchpad_support.sh
/usr/local/etc/chrx-files/etc/skel/.xbindkeysrc
/usr/local/etc/chrx-files/etc/skel/.xsessionrc
/usr/local/etc/chrx-files/etc/sysctl.conf-append
/usr/local/etc/chrx-files/etc/systemd/logind.conf-append
/usr/local/etc/chrx-files/etc/udev/rules.d/10-disable-ehci.rules
/usr/local/etc/chrx-files/etc/udev/rules.d/99-hide-disks-mmc.rules
/usr/local/etc/chrx-files/etc/udev/rules.d/99-hide-disks-nvme.rules
/usr/local/etc/chrx-files/etc/udev/rules.d/99-hide-disks-sda.rules
