# manifest/x86_64/a/atop.filelist
/usr/local/bin/atop
/usr/local/bin/atop-2.10.0
/usr/local/bin/atopcat
/usr/local/bin/atopconvert
/usr/local/bin/atophide
/usr/local/bin/atopsar
/usr/local/bin/atopsar-2.10.0
/usr/local/etc/cron.d/atop
/usr/local/etc/default/atop
/usr/local/etc/init.d/atop
/usr/local/etc/init.d/atopacct
/usr/local/sbin/atopacctd
/usr/local/sbin/atopgpud
/usr/local/share/atop/atop.daily
/usr/local/share/man/man1/atop.1.zst
/usr/local/share/man/man1/atopcat.1.zst
/usr/local/share/man/man1/atopconvert.1.zst
/usr/local/share/man/man1/atophide.1.zst
/usr/local/share/man/man1/atopsar.1.zst
/usr/local/share/man/man5/atoprc.5.zst
/usr/local/share/man/man8/atopacctd.8.zst
/usr/local/share/man/man8/atopgpud.8.zst
/usr/local/var/log/atop/dummy_after
/usr/local/var/log/atop/dummy_before
