# manifest/i686/f/fuse3.filelist
/usr/local/sbin/mount.fuse3
/usr/local/include/fuse3/fuse_opt.h
/usr/local/include/fuse3/fuse_log.h
/usr/local/include/fuse3/fuse_lowlevel.h
/usr/local/include/fuse3/cuse_lowlevel.h
/usr/local/include/fuse3/libfuse_config.h
/usr/local/include/fuse3/fuse.h
/usr/local/include/fuse3/fuse_common.h
/usr/local/share/man/man8/mount.fuse3.8.zst
/usr/local/share/man/man1/fusermount3.1.zst
/usr/local/lib/pkgconfig/fuse3.pc
/usr/local/lib/libfuse3.so.3.14.0
/usr/local/bin/fusermount3
/usr/local/etc/fuse.conf
/usr/local/etc/init.d/fuse3
/usr/local/etc/udev/rules.d/99-fuse3.rules
/usr/local/sbin/fusermount3
/usr/local/lib/libfuse3.so
/usr/local/lib/libfuse3.so.3
