# manifest/x86_64/c/collectl.filelist
/usr/local/bin/collectl
/usr/local/bin/colmux
/usr/local/etc/collectl.conf
/usr/local/etc/init.d/collectl
/usr/local/share/collectl/UNINSTALL
/usr/local/share/collectl/envrules.std
/usr/local/share/collectl/formatit.ph
/usr/local/share/collectl/gexpr.ph
/usr/local/share/collectl/graphite.ph
/usr/local/share/collectl/hello.ph
/usr/local/share/collectl/lexpr.ph
/usr/local/share/collectl/misc.ph
/usr/local/share/collectl/statsd.ph
/usr/local/share/collectl/util/client.pl
/usr/local/share/collectl/vmstat.ph
/usr/local/share/collectl/vmsum.ph
/usr/local/share/collectl/vnet.ph
/usr/local/share/doc/collectl/ARTISTIC
/usr/local/share/doc/collectl/BuddyInfo.html
/usr/local/share/doc/collectl/COPYING
/usr/local/share/doc/collectl/CPUs.html
/usr/local/share/doc/collectl/ColmuxCPU.jpg
/usr/local/share/doc/collectl/ColmuxLustre.jpg
/usr/local/share/doc/collectl/Data-brief.html
/usr/local/share/doc/collectl/Data-detail.html
/usr/local/share/doc/collectl/Data-verbose.html
/usr/local/share/doc/collectl/Data.html
/usr/local/share/doc/collectl/Disks.html
/usr/local/share/doc/collectl/Documentation.html
/usr/local/share/doc/collectl/Environmental.html
/usr/local/share/doc/collectl/Examples.html
/usr/local/share/doc/collectl/Exceptions.html
/usr/local/share/doc/collectl/Export.html
/usr/local/share/doc/collectl/FAQ-collectl.html
/usr/local/share/doc/collectl/FileNaming.html
/usr/local/share/doc/collectl/GPL
/usr/local/share/doc/collectl/Ganglia.jpg
/usr/local/share/doc/collectl/Ganglia2.jpg
/usr/local/share/doc/collectl/Ganglia3.jpg
/usr/local/share/doc/collectl/Gexpr.html
/usr/local/share/doc/collectl/Graphite.html
/usr/local/share/doc/collectl/Hello.html
/usr/local/share/doc/collectl/HiResTime.html
/usr/local/share/doc/collectl/Import.html
/usr/local/share/doc/collectl/Infiniband.html
/usr/local/share/doc/collectl/Inodes.html
/usr/local/share/doc/collectl/InputFiles.html
/usr/local/share/doc/collectl/Interrupts.html
/usr/local/share/doc/collectl/Logging.html
/usr/local/share/doc/collectl/Lustre.html
/usr/local/share/doc/collectl/Matrix.html
/usr/local/share/doc/collectl/Memory.html
/usr/local/share/doc/collectl/Messages.html
/usr/local/share/doc/collectl/Misc.html
/usr/local/share/doc/collectl/Network.html
/usr/local/share/doc/collectl/NfsInfo.html
/usr/local/share/doc/collectl/OpenStack.html
/usr/local/share/doc/collectl/OperationalModes.html
/usr/local/share/doc/collectl/OutputFiles.html
/usr/local/share/doc/collectl/OutputFormats.html
/usr/local/share/doc/collectl/Performance.html
/usr/local/share/doc/collectl/Playback.html
/usr/local/share/doc/collectl/Plotfiles.html
/usr/local/share/doc/collectl/Process.html
/usr/local/share/doc/collectl/RELEASE-collectl
/usr/local/share/doc/collectl/RunningAsAService.html
/usr/local/share/doc/collectl/SlabInfo.html
/usr/local/share/doc/collectl/SlowProc.html
/usr/local/share/doc/collectl/Socket.html
/usr/local/share/doc/collectl/Sockets.html
/usr/local/share/doc/collectl/Startup.html
/usr/local/share/doc/collectl/TheMath.html
/usr/local/share/doc/collectl/Tutorial-Lustre.html
/usr/local/share/doc/collectl/Tutorial.html
/usr/local/share/doc/collectl/WhySummary.html
/usr/local/share/doc/collectl/colmux.html
/usr/local/share/doc/collectl/slub.jpg
/usr/local/share/doc/collectl/style.css
/usr/local/share/man/man1/collectl.1.gz
/usr/local/share/man/man1/colmux.1.gz
